# django、uwsgi、nginx部署纪实/mysite/docker/supervisord.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /data/apps/mysite/docker/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx